# transcription-daemon/geo-transcript-srv.service
[Unit]
Description=OSM l10n transcription server
ConditionPathIsDirectory=/usr/local/share/osml10n/boundaries
 
[Service]
TimeoutStartSec=0
Restart=always
Type=notify
NotifyAccess=all
User=nobody
Group=nogroup
StandardError=null

ExecStart=/usr/local/bin/geo-transcript-srv.py -s -g /usr/local/share/osml10n/boundaries
# Max number of open files
LimitNOFILE=49152

[Install]
WantedBy=multi-user.target
