# transcription-daemon/geo-transcript-srv.service.template
[Unit]
Description=OSM l10n transcription server
 
[Service]
TimeoutStartSec=0
Restart=always
Type=notify
NotifyAccess=all
User=nobody
Group=nogroup

ExecStart=%PYTARGET%/bin/geo-transcript-srv.py -s
# Max number of open files
LimitNOFILE=49152

[Install]
WantedBy=multi-user.target
